# thermal-guard.service
[Unit]
Description=GPU/NVMe thermal guard (auto power-limit stepdown + temp alerts)
After=multi-user.target

[Service]
Type=simple
ExecStartPre=/usr/bin/nvidia-smi -pl 200
ExecStart=/usr/local/bin/thermal-guard.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
